# src/autostart/chrome.desktop
[Desktop Entry]
Type=Application
Name=Google Chrome Custom
Exec= sh -c '$HOME/bin/chrome/google-chrome-stab.sh
Terminal=false

