# src/autostart/chrome.desktop
[Desktop Entry]
Type=Application
Name=Google Chrome Custom
Exec= sh -c '$HOME/bin/google-chrome-stab'
Terminal=false

