# src/autostart/nosessions.desktop
[Desktop Entry]
Type=Application
Name=Clear XFCE Sessions
Exec= sh -c '$HOME/bin/nosessions.sh'
